# infra/notification/php/supervisord/php.conf
[program:php]
command=/usr/sbin/php5-fpm
user = www-data
autostart = true
startsecs=1
startretries=10
stdout_logfile=/var/log/php/php_output.log
stderr_logfile=/var/log/php/php_errors.log
